# infod-client.example.service
[Unit]
Description=Run infod client

[Service]
ExecStart=/usr/sbin/infod_client
Restart=on-failure

[Install]
WantedBy=multi-user.target
