# infod-server.example.service
[Unit]
Description=Run infod server

[Service]
ExecStart=/usr/sbin/infod_server
Restart=on-failure

[Install]
WantedBy=multi-user.target
